# robotVM/build/vrep/home/vagrant/.config/autostart/vrep.desktop
[Desktop Entry]
Name=V-REP Simulator
Exec=/opt/V-REP_PRO_EDU_V3_6_2_Ubuntu18_04/vrep.sh
Type=Application
X-GNOME-Autostart-enabled=true