# robotVM/build/base/home/vagrant/.config/autostart/desktop_settings.desktop
[Desktop Entry]
Name=Apply Gnome Settings
Exec=/opt/set-gnome-settings.sh
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Type=Application
